# files/media_server/etc/containers/systemd/beszel.container
# /etc/containers/systemd/beszel.container
#
# This Quadlet file runs Beszel server monitoring as a systemd service.

[Unit]
Description=Beszel Server Monitoring
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/getbeszel/beszel:latest
ContainerName=beszel

PublishPort=0.0.0.0:8090:8090

# Volume for persistent data
Volume=/var/opt/beszel/data:/beszel_data:Z

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
