# files/media_server/etc/containers/systemd/beszel-agent.container
# /etc/containers/systemd/beszel-agent.container
#
# This Quadlet file runs the Beszel agent as a systemd service.
# The agent collects system and container metrics and sends them to the Beszel hub.

[Unit]
Description=Beszel Monitoring Agent
After=network-online.target beszel.service
Wants=network-online.target

[Container]
ContainerName=beszel-agent
Image=docker.io/henrygd/beszel-agent:latest
Network=host

Environment=LISTEN=45876
Environment=KEY='ssh-ed25519 AAAAC3NzaC1lZDI1NTE5AAAAILDt+Y1uProBNg1aDX+9y0o8jLNwo9LNcIG9fSirfTHr'
Environment=TOKEN=6d38f-80403a16-42322-00e762664
Environment=HUB_URL=http://home.ancon-mimosa.ts.net:8090

# Volume for agent's persistent data
# Ensure the host directory /var/opt/beszel-agent/data exists
Volume=/var/opt/beszel-agent/data:/var/lib/beszel-agent:Z

# Volume mounts for the agent to collect metrics
# Mount Podman socket to read container stats
Volume=/run/podman/podman.sock:/var/run/docker.sock:ro,z
# Mount system directories for host metrics
Volume=/proc:/host/proc:ro,z
Volume=/sys:/host/sys:ro,z

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
