# files/media_server/etc/systemd/system/media-config-subst.service
# /etc/systemd/system/media-config-subst.service
#
# This service sources the configuration script and uses envsubst
# to generate the final config files from templates.
# It runs after tailscale.service is active to ensure the IP is available.

[Unit]
Description=Substitute variables in media config files
# Run after tailscale is up to ensure we can get an IP.
Requires=tailscaled.service
After=tailscaled.service
# Run before the services that need the generated configs.
Before=caddy.service dnsmasq.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/bos/generate-media-configs.sh

[Install]
WantedBy=multi-user.target
