# files/media_server/etc/bos/media_server_templates/caddy.container
# /etc/containers/systemd/caddy.container
#
# This Quadlet file runs the Caddy web server as a systemd service.
# It acts as a reverse proxy for other services.
#
# To use:
# 1. Place this file in /etc/containers/systemd/
# 2. Ensure /etc/caddy/Caddyfile is present on the host.
# 3. Create host directories for persistent data:
#    sudo mkdir -p /opt/caddy/config /opt/caddy/data
# 4. Reload systemd: systemctl daemon-reload
# 5. Start and enable the service: systemctl enable --now caddy.service

[Unit]
Description=Caddy Web Server
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=caddy
Image=docker.io/library/caddy:latest

# Use host networking to allow Caddy to proxy to other containers
# that have published ports on the host.
Network=host

# Volume for Tailscale socket to enable automatic HTTPS for *.ts.net domains
Volume=/var/run/tailscale:/var/run/tailscale:ro,z

# Volumes for Caddyfile, configuration, and data.
# The ":Z" at the end handles SELinux permissions.
Volume=/etc/caddy/Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=/var/opt/caddy/config:/config:Z
Volume=/var/opt/caddy/data:/data:Z

[Service]
Restart=on-failure

[Install]
# Start this service automatically on boot
WantedBy=multi-user.target
