# files/media_server/etc/bos/media_server_templates/dnsmasq.container
# /etc/containers/systemd/dnsmasq.container
#
# This Quadlet file runs the dnsmasq DNS server.
#
# To use:
# 1. Edit /etc/dnsmasq.d/hosts.conf with your server's Tailscale IP.
# 2. Reload systemd: systemctl daemon-reload
# 3. Start and enable the service: systemctl enable --now dnsmasq.service

[Unit]
Description=Dnsmasq DNS Server
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/jpillora/dnsmasq:latest
ContainerName=dnsmasq
User=dnsmasq
AddCapability=CAP_NET_BIND_SERVICE

# Publish DNS port 53 on both TCP and UDP, but only on the Tailscale IP.
# --- IMPORTANT ---
# Replace "100.X.Y.Z" with your server's Tailscale IP address (`tailscale ip -4`)
PublishPort=100.X.Y.Z:53:53/tcp
PublishPort=100.X.Y.Z:53:53/udp

# Mount the configuration file. The default config in the image
# is set up to read from /etc/dnsmasq.d/
Volume=/opt/dnsmasq/dnsmasq.d:/etc/dnsmasq.d:Z

Restart=on-failure

[Install]
WantedBy=multi-user.target
