# files/media_server/etc/bos/media_server_templates/homepage.container
# /etc/containers/systemd/homepage.container
#
# This Quadlet file runs Homepage as a systemd service.

[Unit]
Description=Homepage Dashboard
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/gethomepage/homepage:latest
ContainerName=homepage
PublishPort=127.0.0.1:3000:3000

# Volume for persistent configuration
Volume=/var/opt/homepage/config:/app/config:Z

# Volume for podman socket to get container status
Volume=/var/run/podman/podman.sock:/var/run/docker.sock:Z

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target
